# configs/programs/lab.conf
# ~ supervisord config ~
# Change the file name to `disabled-...` to disable the program below

[program:lab]
directory=/lab
command=/bin/bash -c "node dist-backend/server.js"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
